# systemd/orion-dns-backup.service
[Unit]
Description=Orion DNS Backup Service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/orion-dns-ha
ExecStart=/opt/orion-dns-ha/ops/orion-dns-backup.sh
