# systemd/orion-dns-sync.service
[Unit]
Description=Orion DNS Pi-hole Sync Service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/orion-dns-ha
ExecStart=/opt/orion-dns-ha/ops/pihole-sync.sh
