# systemd/orion-dns-ha-sync.timer
[Unit]
Description=Orion DNS HA - Pi-hole Config Sync Timer
Documentation=https://github.com/orionsentinel/Orion-sentinel-ha-dns
Requires=orion-dns-ha-sync.service

[Timer]
# Run sync every hour
OnCalendar=hourly
# Run 5 minutes after boot
OnBootSec=5min
# Run immediately if we missed a scheduled run
Persistent=true

[Install]
WantedBy=timers.target
