# systemd/pihole-sync.timer
[Unit]
Description=Pi-hole Configuration Sync Timer
Requires=pihole-sync.service

[Timer]
# Run every 6 hours
OnCalendar=*-*-* 00,06,12,18:00:00
# If the system was down during a scheduled run, run the service once after boot
Persistent=true
# Add a random delay of up to 5 minutes to avoid synchronized load
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
