# systemd/pihole-sync.service
[Unit]
Description=Pi-hole Configuration Sync Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/orion-dns-ha/Orion-sentinel-ha-dns
EnvironmentFile=/opt/orion-dns-ha/Orion-sentinel-ha-dns/.env
ExecStart=/opt/orion-dns-ha/Orion-sentinel-ha-dns/scripts/pihole-sync.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
